# ansible/roles/showroom/templates/showroom.service.j2
[Unit]
Description=Showroom Service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
User={{ showroom_user | default('showroom') }}
Group={{ showroom_user_group | default('showroom') }}
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=COMPOSE_FILE={{ showroom_user_home_dir }}/orchestration/container-compose.yml
Environment=SECURITY_OPT="--security-opt label=type:container_runtime_t"
Environment=CLEANUP_OPT="--remove-orphan"
ExecStart=/usr/local/bin/podman-compose -f $COMPOSE_FILE --podman-run-args $SECURITY_OPT up -d
ExecStop=/usr/local/bin/podman-compose -f $COMPOSE_FILE down
Restart=on-failure
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
