# ansible/roles_ocp_workloads/ocp4_workload_vscode_bastion/templates/code-server.service.j2
[Unit]
Description=Code server system service for user {{ ocp4_workload_vscode_bastion_user_name }}
After=network.target

[Service]
User={{ ocp4_workload_vscode_bastion_user_name }}
Group={{ ocp4_workload_vscode_bastion_user_group }}

Type=exec
ExecStart=/usr/bin/code-server
Restart=always

[Install]
WantedBy=multi-user.target
