# ansible/roles_ocp_workloads/ocp4_workload_mad_roadshow_v2/files/oracle/update_ipaddress.service
[Unit]
Description=Run script at startup after network becomes reachable
After=network.target
Before=oracle-xe-21c.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/update_ipaddress.sh
TimeoutStartSec=0

[Install]
WantedBy=default.target
