# ansible/roles/config-redis/templates/redis.service.j2
[Unit]
Description=Redis Docker Container
After=docker.service
Requires=docker.service

[Service]
EnvironmentFile={{ systemd_environmentfile_dir }}/{{ redis_name }}
TimeoutStartSec=0
Restart=always
ExecStartPre=/usr/bin/docker pull ${IMAGE}
ExecStart=/usr/bin/docker run --rm -v ${REDIS_HOST_STORAGE_DIR}:${REDIS_CONTAINER_STORAGE_DIR}:Z -p ${REDIS_CONTAINER_PORT}:${REDIS_HOST_PORT} --name %n ${IMAGE}
ExecStop=/usr/bin/docker stop %n
ExecReload=/usr/bin/docker restart %n

[Install]
WantedBy=multi-user.target