# ansible/roles/config-quay-builder/templates/quay-builder.service.j2
[Unit]
Description=Quay Builder Docker Container
After=docker.service
Requires=docker.service

[Service]
EnvironmentFile={{ systemd_environmentfile_dir }}/{{ quay_builder_name }}
TimeoutStartSec=0
Restart=always
ExecStartPre=/usr/bin/docker pull ${QUAY_BUILDER_IMAGE}
ExecStart=/usr/bin/docker run --rm --name %n -e SERVER={% if quay_builder_ssl_trust_configure|bool %}wss{% else %}ws{% endif %}://${QUAY_ADDRESS} {% if quay_builder_ssl_trust_configure|bool %} -v ${QUAY_BUILDER_SSL_TRUST_HOST_FILE}:${QUAY_BUILDER_SSL_TRUST_CONTAINER_FILE}:Z {% endif %} -v /var/run/docker.sock:/var/run/docker.sock --entrypoint=/bin/sh ${QUAY_BUILDER_IMAGE} -c '{% if quay_builder_ssl_trust_configure|bool %}/usr/sbin/update-ca-certificates && {% endif %}quay-builder'
ExecStop=/usr/bin/docker stop %n
ExecReload=/usr/bin/docker restart %n

[Install]
WantedBy=multi-user.target