# ansible/configs/ansible-advanced-v2/files/all-in-one-pod.service.j2
[Unit]
Description=Podman all-in-one-pod Service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/usr/local/bin/podman-compose -f /opt/podman/resources/orchestration/all-in-one-podman-compose.yml up -d
ExecStop=/usr/local/bin/podman-compose -f /opt/podman/resources/orchestration/all-in-one-podman-compose.yml down
KillMode=none
Type=forking
TimeoutSec=900

[Install]
WantedBy=multi-user.target default.target
